# scheduled/telegram-triage/evening.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sheldon.telegram-triage.evening</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>__RUNNER__</string>
    <string>evening</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>19</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/sheldon-telegram-triage/evening.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/sheldon-telegram-triage/evening.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>__REPO__</string>

  <key>AbandonProcessGroup</key>
  <true/>
</dict>
</plist>
