# scheduled/telegram-triage/morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sheldon.telegram-triage.morning</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>__RUNNER__</string>
    <string>morning</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <!-- launchd's default: if the Mac is asleep/off at fire time, the job runs
       as soon as it wakes. No StartInterval needed. -->
  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/sheldon-telegram-triage/morning.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/sheldon-telegram-triage/morning.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>__HOME__</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>__REPO__</string>

  <!-- Launch even if the scheduled time was missed while the Mac was off -->
  <key>AbandonProcessGroup</key>
  <true/>
</dict>
</plist>
